# config/systemd/spotify-notify/spotify-notify.service
[Unit]
Description=Spotify Track Change Notifications

[Service]
ExecStart=%h/.local/bin/spotify-notify.sh

[Install]
WantedBy=default.target

